# ops/llm-health-check.service
[Unit]
Description=LLM Benchmark Daily Health Check
After=network-online.target
Wants=network-online.target
OnFailure=llm-health-check-failure@%n.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/llm-benchmarks

# Load environment (MongoDB URI, OpenAI key, etc.)
EnvironmentFile=/etc/llm-benchmarks/health-check.env

# Run the health check script
ExecStart=/usr/bin/python3 /opt/llm-benchmarks/ops/daily-health-check.py --days 1

# Resource limits
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Timeout after 5 minutes
TimeoutStartSec=300
