# ops/llm-discovery.service
[Unit]
Description=LLM Model Discovery via OpenRouter
After=network-online.target
Wants=network-online.target
OnFailure=llm-health-check-failure@%n.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/llm-benchmarks

# Load environment (MongoDB URI, OpenAI key, etc.)
EnvironmentFile=/etc/llm-benchmarks/health-check.env

# Run discovery fetch
ExecStart=/usr/bin/python3 -m api.llm_bench.discovery.cli fetch

# Resource limits
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Timeout after 5 minutes
TimeoutStartSec=300
